# package/rooter/0optionalapps/ext-domain/files/etc/init.d/domain
#!/bin/sh /etc/rc.common

START=88

start_service()
{
	/usr/lib/domain/filter.sh
}