# package/rooter/0optionalapps/ext-domain/files/etc/init.d/domain
#!/bin/sh /etc/rc.common

START=91

log() {
	logger -t "Domain Filter " "$@"
}

start()
{
	#log "Start Domain"
}