# package/rooter-extra/gl-uart/files/etc/init.d/x1200_init
#!/bin/sh /etc/rc.common

START=10

start(){
	uartinit 9600 /dev/ttyUSB0
	uartinit 9600 /dev/ttyUSB1
}


