# configfiles/huasifei/files1/etc/uci-defaults/03-tzone
#!/bin/sh

uci set system.@system[0].hostname='Huasifei'
uci commit system
/etc/init.d/system restart