# package/rooter/ext-rooter-basic/files/etc/uci-defaults/30-dhcpmaster
#!/bin/sh

uci set dhcp.lan.master='1'
uci set dhcp.lan.ra='hybrid'
uci set dhcp.lan.ra_flags='none'
uci set dhcp.lan.dhcpv6='hybrid'
uci commit dhcp
/etc/init.d/dnsmasq restart
