# configfiles/Mesh/files-master/etc/uci-defaults/03-tzone
#!/bin/sh

uci set system.@system[0].hostname='ROOter-Mesh-Wifi'
uci commit system
/etc/init.d/system restart
