# package/rooter/0optionalapps/ext-zerotier/files/etc/uci-defaults/64-zerotier
#!/bin/sh

id=$(uci -q get custom.zerotier.networkid)
uci delete zerotier.zerotier1.join
uci add_list zerotier.zerotier1.join=$id
uci commit zerotier
/etc/init.d/zerotier restart
