# feeds/packages/net/zerotier/files/etc/uci-defaults/80-zt-migration
#!/bin/sh

OLD="zerotier"

nets=$(uci -q get $OLD.zerotier.join)

if [ -n "$nets" ]; then
	enable=$(uci -q get $OLD.zerotier.enabled)
	secret=$(uci -q get $OLD.zerotier.secret)
	
	uci delete $OLD.zerotier
	uci set $OLD.global=zerotier
	uci set $OLD.global.enabled=$enable
	uci set $OLD.global.secret=$secret
	
	uci set $OLD.zerotier=network
	uci set $OLD.zerotier.id=$nets
	uci commit $OLD
fi
